# proxy_gen_lib/squidsetup.service
[Unit]
Description=Setup squid

[Service]
Type=simple
ExecStart=/usr/bin/proxystartupscript

[Install]
WantedBy=multi-user.target